# DB VM/fastapi-mongodb-app.conf
[program:fastapi-mongodb-app]

command=/home/user/employee_db/gunicorn_start
user=user
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/home/user/mongodb/logs/gunicorn-error.log
